# crontab/crontab_job
* * * * * curl -X GET http://localhost:5000/motion/auto-activate >/dev/null 2>&1
